# deploy/podman/keyspiderd.container
# Podman Quadlet (systemd-managed container)
# Place as: /etc/containers/systemd/keyspiderd.container

[Unit]
Description=keyspiderd (podman)
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/keyspider:latest
ContainerName=keyspiderd
PublishPort=8080:8080
Environment=KEYSPIDER_DB_DSN=postgres://keyspider:keyspider@host.containers.internal:5432/keyspider?sslmode=disable
# Mount SSH config and (optionally) your SSH keys.
Volume=/etc/ssh/ssh_config:/etc/ssh/ssh_config:ro
# Example (adjust):
# Volume=%h/.ssh:/root/.ssh:ro

[Install]
WantedBy=default.target
